# scripts/com.mohltbot.moonshot-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mohltbot.moonshot-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mohlt/.openclaw/workspace/scripts/sync-moonshot-balance.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/mohlt/.openclaw/workspace/logs/moonshot-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mohlt/.openclaw/workspace/logs/moonshot-sync-error.log</string>
</dict>
</plist>
